# .docker/cron/crontab
# Delete all expired notes from the database every 30 minutes
*/30 * * * * /app/scripts/vacuum-notes
